# sample_configuration/discoroot/etc/cron.d/cloudwatch
*/1 * * * * root /opt/wgen/bin/upload_cloudwatch.sh &> /dev/null
