# sample_configuration/discoroot/etc/cron.d/internal_smoketest
*/5 * * * * root /opt/wgen/bin/smokealarm.sh 2>&1 | /usr/bin/logger -t disco.smoketest -p local0.info
