# sample_configuration/discoroot/etc/cron.d/ipsec_status
PATH=/bin:/usr/bin:/usr/sbin
*/5 * * * * root /opt/wgen/bin/ipsec_status.sh &> /dev/null
