# sample_configuration/discoroot/etc/cron.d/monitor_processes
*/15 * * * * root /opt/wgen/bin/monitor_processes.sh &> /dev/null
