# sample_configuration/discoroot/etc/cron.d/monitor_connectivity
*/2 * * * * root /opt/wgen/bin/monitor_connectivity.sh &> /dev/null
